# docker/api/service/api/run
#!/bin/sh
cd /home/app

exec /sbin/setuser app passenger start --port $WEB_PORT --min-instances $PASSENGER_MIN_INSTANCES
