# docker/api/service/api/run
#!/bin/sh
cd /home/app

if [ "$RAILS_ENV" = "development" ]; then
    exec /sbin/setuser app passenger start --port $WEB_PORT --min-instances $PASSENGER_MIN_INSTANCES
else
    exec /sbin/setuser app passenger start --port $WEB_PORT --min-instances $PASSENGER_MIN_INSTANCES --no-friendly-error-pages
fi